# upstart/agentbox-chief.conf
#
# install with:
# sudo initctl start agentbox-chief
#
start on runlevel [2345]
# setuid agentbox
chdir /home/agentbox/chief
exec ./chief -p 2010
respawn
respawn limit 10 90

