# upstart/agentbox-control.conf
#
# install with:
# sudo initctl start agentbox-control
#
start on runlevel [2345]
# setuid agentbox
chdir /home/agentbox/control
exec ./control -p 2010
respawn
respawn limit 10 90

